# systemFiles/supervisor/conf.d/blueberry-celerybeat.conf
[program:blueberry-celerybeat]
; Set full path to celery program if using virtualenv
command=celery beat -A webAPI --loglevel=DEBUG

; remove the -A myapp argument if you are not using an app instance

directory=/var/www/BlueberryPie/app
user=pi
numprocs=1
stdout_logfile=/var/log/raspberry/celery-beat.log
stderr_logfile=/var/log/raspberry/celery-beat.log
autostart=true
autorestart=true
startsecs=10

; if rabbitmq is supervised, set its priority higher
; so it starts first
priority=999
